# sonacip.service
[Unit]
Description=SONACIP Flask Application
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/sonacip
Environment=PATH=/root/sonacip/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/root/sonacip
Environment=HOME=/root

# Gunicorn with wsgi entrypoint
ExecStart=/root/sonacip/venv/bin/gunicorn -c /root/sonacip/gunicorn.conf.py wsgi:app

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true

# Restart policy - stable, no rapid restarts
Restart=on-failure
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
