# deploy/sonacip.service
[Unit]
Description=SONACIP - Piattaforma Gestione Società Sportive
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=sonacip
Group=sonacip
WorkingDirectory=/var/www/sonacip
EnvironmentFile=/var/www/sonacip/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/sonacip/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 120 \
    --access-logfile /var/www/sonacip/logs/access.log \
    --error-logfile /var/www/sonacip/logs/error.log \
    --capture-output \
    wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=10

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/sonacip
UMask=0027

[Install]
WantedBy=multi-user.target
