# deploy/sonacip-calendar-reminders.timer
[Unit]
Description=Run SONACIP calendar reminders periodically

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target

