# deploy/sonacip-compliance.service
[Unit]
Description=SONACIP Compliance (certificates/fees) Job
After=network.target

[Service]
Type=oneshot
User=sonacip
Group=sonacip
WorkingDirectory=/opt/sonacip
EnvironmentFile=/opt/sonacip/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sonacip/venv/bin/python /opt/sonacip/scripts/daily_compliance_job.py

