# deploy/sonacip.service
[Unit]
Description=SONACIP Gunicorn Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/sonacip
EnvironmentFile=-/opt/sonacip/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sonacip/venv/bin/gunicorn -c /opt/sonacip/gunicorn.conf.py run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
