# deploy/sonacip-maintenance.service
[Unit]
Description=SONACIP Maintenance Job (calendar+compliance)
After=network.target

[Service]
Type=oneshot
User=sonacip
Group=sonacip
WorkingDirectory=/opt/sonacip
EnvironmentFile=/opt/sonacip/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sonacip/venv/bin/python /opt/sonacip/scripts/maintenance_job.py

