# deploy/sonacip-backup.timer
[Unit]
Description=Run SONACIP backups daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
