# deployment/sonacip.service
[Unit]
Description=Gunicorn instance to serve SONACIP
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/sonacip
Environment="PATH=/opt/sonacip/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/opt/sonacip/venv/bin/gunicorn -c gunicorn_config.py run:app

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:/opt/sonacip/logs/sonacip.log
StandardError=append:/opt/sonacip/logs/sonacip_error.log

[Install]
WantedBy=multi-user.target
